# my-bot/bot/crontab
# ┌───────── minute
# │ ┌─────── hour (24h)
# │ │ ┌───── day of month
# │ │ │ ┌─── month
# │ │ │ │ ┌─ day of week (0=Sunday or SUN)
# │ │ │ │ │
# │ │ │ │ │
# Run every 10 min from 12:10..12:59 (Sun..Wed)
10-59/10 12 * * SUN,MON,TUE,WED  python /app/book_appointment_real.py >> /proc/1/fd/1 2>&1
# Run every 10 min from 13:00..13:59 (Sun..Wed)
*/10      13 * * SUN,MON,TUE,WED  python /app/book_appointment_real.py >> /proc/1/fd/1 2>&1
# Run once at 14:00 (Sun..Wed)
0         14 * * SUN,MON,TUE,WED  python /app/book_appointment_real.py >> /proc/1/fd/1 2>&1